# systemd/boot-report.timer
[Unit]
Description=Run Boot server observability report every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=30s
Persistent=true
Unit=boot-report.service

[Install]
WantedBy=timers.target
