# systemd/boot-report.service
[Unit]
Description=Telegram Boot Report
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
User=bootreport
Group=bootreport
ExecStart=/opt/boot-report/boot-report.sh
TimeoutStartSec=240

# Runtime dir en /run con permisos correctos para user no-root
RuntimeDirectory=boot-report
RuntimeDirectoryMode=0750

# Reintento si falla Telegram/red
Restart=on-failure
RestartSec=30

# Hardening (baseline compatible con curl + lectura de sistema)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
UMask=0077
SystemCallArchitectures=native
SystemCallFilter=@system-service @network-io

[Install]
WantedBy=multi-user.target
