# systemd/boot-report.service
[Unit]
Description=Send daily report to Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=bootreport
Group=bootreport

# /var/lib/boot-report (persistente) para stamp diario
StateDirectory=boot-report
StateDirectoryMode=0750

# /run/boot-report para lock temporal
RuntimeDirectory=boot-report
RuntimeDirectoryMode=0750
Environment=XDG_RUNTIME_DIR=/run/boot-report

ExecStart=/opt/boot-report/boot-report.sh
TimeoutStartSec=240

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true

ReadOnlyPaths=/proc /sys
ReadWritePaths=/run/boot-report /var/lib/boot-report

[Install]
WantedBy=multi-user.target
