# systemd/boot-report.service
[Unit]
Description=Boot server observability report
Documentation=file:/opt/boot-report/README-INSTALL.md
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/opt/boot-report/config/server.env
EnvironmentFile=-/etc/boot-report/boot-report.env
ExecStart=/opt/boot-report/bin/boot-report
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/boot-report /opt/boot-report
LockPersonality=true
MemoryDenyWriteExecute=true
